# platform/runtime/web/backend/unibos.service
[Unit]
Description=UNIBOS Django Backend
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=gunicorn
WorkingDirectory=/var/www/unibos/backend
EnvironmentFile=/var/www/unibos/backend/.env
ExecStart=/var/www/unibos/backend/venv/bin/gunicorn unibos_backend.asgi:application --config /var/www/unibos/backend/gunicorn.conf.py
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/unibos/backend/media /var/www/unibos/backend/logs /var/log/unibos
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="DJANGO_SETTINGS_MODULE=unibos_backend.settings.production"

[Install]
WantedBy=multi-user.target